# chef/site-cookbooks/landcare/templates/default/upstart-uwsgi.conf.erb
description "uWSGI server for <%= node[:app][:name] %>" 

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec uwsgi --ini <%= node[:app][:paths][:site] %>/uwsgi.ini
